# templates/default/horizon-25/opennms.service.erb
[Unit]
Description=OpenNMS server
Wants=postgresql.service
Requires=network.target network-online.target
After=postgresql.service postgresql-10.service postgresql-11.service postgresql-12.service network.target network-online.target

[Service]
User=root
Environment="OPENNMS_HOME=/opt/opennms"

Type=simple
PIDFile=/opt/opennms/logs/opennms.pid

ExecStart=/etc/init.d/opennms -f <%=@start_opts%> start
ExecStop=/etc/init.d/opennms stop
TimeoutStartSec=<%=@timeout_start_sec%>

[Install]
WantedBy=multi-user.target
