# templates/default/horizon-16/opennms.service.erb
[Unit]
Description=OpenNMS server
Wants=postgresql.service
Requires=network.target network-online.target
After=postgresql.service network.target network-online.target

[Service]
User=root

Type=forking
PIDFile=/opt/opennms/logs/opennms.pid

ExecStart=/opt/opennms/bin/opennms <%=@start_opts%> start
ExecStop=/opt/opennms/bin/opennms stop
<% if @timeout_start_sec != '' -%>
TimeoutStartSec=<%=@timeout_start_sec%>
<% end -%>

[Install]
WantedBy=multi-user.target
